# supervisord.conf
[supervisord]
nodaemon=true

[program:avahi]
command=avahi-daemon --no-rlimits

[program:timemachine]
command=/start_netatalk.sh